# supervisor/conf.d/uwsgi.conf
[program:uwsgi]
command = /usr/local/bin/uwsgi --http 127.0.0.1:40000
          --wsgi-file /opt/uwsgi/sample.py
          --master --processes 1 --threads 2
user=root
stopsignal=INT
