# supervisor/conf.d/uwsgi.conf
[program:uwsgi]
command = /usr/local/bin/uwsgi --http 0.0.0.0:40000
          --wsgi-file %(ENV_WSGI_FILE)s
          --master --processes 1 --threads 2
user=root
stopsignal=INT
